# ansible/roles/myrs/files/myrs.service
[Unit]
Description=myrs hobby rental server implements

[Service]
ExecStart=/usr/local/bin/myrs

[Install]
WantedBy=multi-user.target