# ansible/roles/myrs/templates/myrs.service.j2
[Unit]
Description=myrs hobby rental server implements

[Service]
ExecStart=/usr/local/bin/myrs -site-domain {{ site_domain }}

[Install]
WantedBy=multi-user.target